# packaging/deb-systemd/debian/mackerel-agent.prerm
#!/bin/sh
set -e

case "$1" in
remove)
  # TODO how to stop systemd daemon
  # if [ -x "/etc/init.d/mackerel-agent" ]; then
  #   invoke-rc.d mackerel-agent stop || exit $?
  # fi
;;
upgrade|deconfigure|failed-upgrade)
  exit 0
;;
*)
  echo "postinst called with unknown argument \`$1'" >&2
  exit 1
;;
esac

#DEBHELPER#
