# packaging/deb/debian/mackerel-agent.prerm
#!/bin/sh
set -e
if [ -x "/etc/init.d/mackerel-agent" ]; then
  invoke-rc.d mackerel-agent stop || exit $?
fi

#DEBHELPER#
